# supervisor/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; path to your socket file

[supervisord]
logfile=logs/supervisord.log ; supervisord log file
logfile_maxbytes=50MB       ; maximum size of logfile before rotation
logfile_backups=10          ; number of backed up logfiles
loglevel=info               ; info, debug, warn, trace
pidfile=logs/supervisord.pid ; pidfile location
nodaemon=false              ; run supervisord as a daemon
minfds=1024                 ; number of startup file descriptors
minprocs=200                ; number of process descriptors
user=penny                   ; default user
childlogdir=logs           ; where child log files will live


[program:celeryd]
command=venv/bin/celery -A troupon worker -l debug
stdout_logfile=logs/celeryd.log
stderr_logfile=logs/celeryd.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600


[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use unix:// schem for a unix sockets.


[include]

files=celeryd.conf
